# php-supervisor/filesystem/etc/supervisor.d/cron.ini
[program:cron]
command=crond -f
autostart=%(ENV_SUPERVISOR_ENABLE_CRON)s
stdout_events_enabled=true
stderr_events_enabled=true
